# dummy.service
[Unit]
Description=Dummy Service
After=network-online.target
OnFailure=notify-telegram@%i.service

[Service]
...
